# rc.local
#!/bin/bash

FILE_PATH="/first_run.sh"
if [ -f "$FILE_PATH" ]; then
    chmod +x "$FILE_PATH"
    ./"$FILE_PATH"
else
    echo "rc.local - Not a first run"
fi

raspi-config nonint do_wifi_country GB
ifconfig wlan0 up

nmcli device set wlan0 managed yes

sysctl -w net.ipv4.ip_forward=1


nmcli device set usb0 managed yes
nmcli con delete usb0-share
nmcli con add type ethernet ifname usb0 con-name usb0-share ip4 192.168.8.1/24
nmcli con mod usb0-share ipv4.method shared
nmcli con up usb0-share
